# core/files/supervisor-program.conf
[program:{{ program_name }}]
command={{ command }}
directory={{ directory }}
autostart=true
autorestart=true
redirect_stderr=true
user={{ user }}
